# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_events_enabled=true
stderr_events_enabled=true

[program:gunicorn]
command=gunicorn --bind 127.0.0.1:8080 --access-logfile /var/log/guni_access.log --error-logfile /var/log/guni_error.log app:callback
directory=/opt/linebot
